# gui_and_audio/rc.local
#!/bin/sh
/home/root/audio/play.init &
exit 0
