# audio/rc.local
#!/bin/bash
/home/root/play.init &
exit 0
