# target-arm-scripts/rc.local
#!/bin/sh

/home/root/run.all &

exit 0
